# Raspbbery_pi_os_config/dash.service
[Unit]
Description=Kivy Dashboard
After=network.target
#After=multi-user.target

[Service]
ExecStartPre=/usr/local/bin/myscript.sh # scpript before staring the app
ExecStart=/usr/bin/python3 /home/dash/dash/main.py
WorkingDirectory=/home/dash
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi
Environment="PYTHONUNBUFFERED=1"
Environment="SDL_VIDEODRIVER=kmsdrm"
Environment="KIVY_GL_BACKEND=kivy.core.window.window_sdl2"

[Install]
WantedBy=graphical.target