# Raspbbery_pi_os_config/dash.service
[Unit]
Description=Kivy Dashboard
After=startcan.service

[Service]
Type=simple
User=dash
Group=dash
WorkingDirectory=/home/dash
ExecStart=/usr/bin/python3 /home/dash/dash/main.py
StandardOutput=inherit
StandardError=inherit
Restart=always

[Install]
WantedBy=multi-user.target